# edubot/supervisord.conf
[program:edubot]
user=root
directory=/home/edubot
command=/home/edubot/edubot.py
autostart=true
autorestart=true
stderr_logfile=/var/log/edubot.err.log
stdout_logfile=/var/log/edubot.out.log
